# debian/nyancat-server.socket
[Unit]
Description=nyancat-server socket

[Socket]
ListenStream=23
Accept=true

[Install]
WantedBy=sockets.target
